# build/profiles/default/packages/middleware/filelist
/usr/local/lib/datastore
/usr/local/lib/dispatcher
/usr/local/lib/freenascli
/usr/local/lib/etcd
/usr/local/lib/networkd
/usr/local/lib/fnstatd
/usr/local/lib/python2.7/site-packages/datastore
/usr/local/lib/python2.7/site-packages/dispatcher
/usr/local/lib/python2.7/site-packages/libzfs
/usr/local/lib/python2.7/site-packages/netif
/usr/local/lib/python2.7/site-packages/cam
/usr/local/lib/python2.7/site-packages/fnutils
/usr/local/lib/python2.7/site-packages/sysctl
/usr/local/lib/python2.7/site-packages/libzfs
/usr/local/lib/python2.7/site-packages/pySMART
/usr/local/share/datastore
/usr/local/etc/rc.d/datastore
/usr/local/etc/rc.d/dispatcher
/usr/local/etc/rc.d/etcd
/usr/local/etc/rc.d/networkd
/usr/local/etc/rc.d/fnstatd
/usr/local/etc/middleware.conf
/usr/local/bin/cli
/usr/local/bin/logincli
/usr/local/sbin/dsutil
/usr/local/sbin/dsdump
/usr/local/sbin/dsrestore
/usr/local/sbin/dspasswd
/usr/local/sbin/dsmigrate
/usr/local/sbin/dsinit
/usr/local/sbin/dispatcherd
/usr/local/sbin/dispatcherctl
/usr/local/sbin/etcd
/usr/local/sbin/fnstatd
/usr/local/sbin/networkd