# build/profiles/freenas10/packages/middleware/filelist
/usr/local/lib/datastore
/usr/local/lib/dispatcher
/usr/local/lib/dsd
/usr/local/lib/etcd
/usr/local/lib/migrate93
/usr/local/lib/networkd
/usr/local/lib/fnstatd
/usr/local/lib/crashd
/usr/local/lib/python3.4/site-packages/datastore
/usr/local/lib/python3.4/site-packages/libzfs
/usr/local/lib/python3.4/site-packages/netif
/usr/local/lib/python3.4/site-packages/cam
/usr/local/lib/python3.4/site-packages/ctl
/usr/local/lib/python3.4/site-packages/cli
/usr/local/lib/python3.4/site-packages/freenas
/usr/local/lib/python3.4/site-packages/sysctl
/usr/local/lib/python3.4/site-packages/libzfs
/usr/local/lib/python3.4/site-packages/pySMART
/usr/local/lib/python3.4/site-packages/ws4py
/usr/local/lib/python3.4/site-packages/ipfs-api
/usr/local/lib/node_modules/dispatcher-client
/usr/local/share/datastore
/usr/local/etc/rc.d/datastore
/usr/local/etc/rc.d/dispatcher
/usr/local/etc/rc.d/dsd
/usr/local/etc/rc.d/etcd
/usr/local/etc/rc.d/networkd
/usr/local/etc/rc.d/fnstatd
/usr/local/etc/middleware.conf
/usr/local/bin/cli
/usr/local/bin/logincli
/usr/local/sbin/crashd
/usr/local/sbin/crash-wrapper
/usr/local/sbin/dsutil
/usr/local/sbin/dsdump
/usr/local/sbin/dsprinttask
/usr/local/sbin/dsrestore
/usr/local/sbin/dspasswd
/usr/local/sbin/dsmigrate
/usr/local/sbin/dsinit
/usr/local/sbin/dispatcherd
/usr/local/sbin/dispatcherctl
/usr/local/sbin/sdsinit
/usr/local/sbin/dsd
/usr/local/sbin/etcd
/usr/local/sbin/fnstatd
/usr/local/sbin/networkd
/usr/local/sbin/migrate93
/usr/local/libexec/taskproxy
